# script/minerProxy.conf

# https://github.com/minerproxyvip/mp.git
[program:minerProxy]
command=/root/minerProxy/minerProxy
directory=/root/minerProxy
autostart=true
autorestart=true






