# script/mp.conf

# https://github.com/minerproxyvip/mp.git
[program:mp]
command=/root/mp/mp
directory=/root/mp
autostart=true
autorestart=true






